# scripts/afd.target
[Unit]
Description=Automatic File Distributor (AFD)

[Install]
WantedBy=multi-user.target
