# scripts/afdmon.target
[Unit]
Description=Monitor for Automatic File Distributor (AFD)

[Install]
WantedBy=multi-user.target
